# config/deployment/student-records-app.service
[Unit]
Description=Student Records App
After=network.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/var/www/student_records_app
Environment=RAILS_ENV=production
Environment=BUNDLE_GEMFILE=/var/www/student_records_app/Gemfile
ExecStart=/usr/local/bin/bundle exec puma -C config/puma.rb
ExecReload=/bin/kill -USR1 $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/www/student_records_app/tmp /var/www/student_records_app/log /var/www/student_records_app/storage
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
MemoryMax=1G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
